# .config/systemd/user/emacs.service.d/override.conf
[Service]
Environment="DOOMDIR=/home/uroborosikes/.config/doom"
Environment="DOOMLOCALDIR=/home/uroborosikes/.emacs.d/local/"